# config/systemd/qcal-automator.service
[Unit]
Description=QCAL piCODE Sovereign Swarm Automator and DLC Bridge
After=network.target
[Service]
Type=simple
User=root
WorkingDirectory=/root/repo_P-NP
ExecStart=/usr/bin/python3 /root/repo_P-NP/scripts/a2a_dlc_automator.py
Restart=always
RestartSec=30
StandardOutput=append:/var/log/qcal/automator.log
StandardError=append:/var/log/qcal/automator.error.log
[Install]
WantedBy=multi-user.target
